# yodad.service
[Unit]
Description=aaaargh

[Service]
ExecStart=yoda
Type=simple

[Install]
WantedBy=default.target
