# yodad.service
[Unit]
Description=aaaargh
After=multi-user.target

[Service]
ExecStart=yoda

[Install]
WantedBy=multi-user.target